# com.projectignite.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectignite.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jw/Projects/Project-Ignite/mcp_main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jw/Projects/Project-Ignite</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jw/Projects/Project-Ignite/mcp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jw/Projects/Project-Ignite/mcp.log</string>
</dict>
</plist> 